# data/hud.service.in
[Unit]
Description=Unity HUD

[Service]
Type=dbus
BusName=com.canonical.hud
ExecStart=@pkglibexecdir@/hud-service
