# data/window-stack-bridge.conf.in
description "Window Stack Bridge"

start on started dbus
stop on desktop-end

exec @pkglibexecdir@/window-stack-bridge
